# src/etc/cron.daily/open.files.sh
#!/bin/sh

/usr/local/emhttp/plugins/open.files/daily.php 1>/dev/null 2>&1
